# zennlogic_ai/infra/systemd/genai-api.service
[Unit]
Description=zennlogic_ai_service REST API
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/env bash -lc 'uv run uvicorn src.zennlogic_ai_service.rest.app:app --host 0.0.0.0 --port 8000 --workers 1 --loop uvloop --http httptools --timeout-keep-alive 15'
Restart=always
RestartSec=5
User=ec2-user
WorkingDirectory=/home/ec2-user/zennlogic_ai_service
EnvironmentFile=/home/ec2-user/zennlogic_ai_service/.env

[Install]
WantedBy=multi-user.target
